# config/udev/99-printer.rules
# udev rules for USB printer hotplug detection
#
# IMPORTANT: Only ONE rule restarts CUPS per hotplug event. The generic
# bInterfaceClass=07 rule catches ALL USB printers (Brother, HP, Canon,
# Epson, etc.) so vendor-specific rules are not needed and were removed
# to prevent CUPS from restarting 5-6 times on a single plug-in event.

# --- Printer connection: single consolidated rule ---
# Matches any USB printer (class 07), restarts CUPS once, regenerates
# Avahi services, disables USB autosuspend, and logs the event.
SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="07", ACTION=="add", \
    RUN+="/bin/systemctl restart cups.service", \
    RUN+="/usr/bin/logger -t printer-hotplug 'USB printer connected, restarting CUPS'", \
    MODE="0666"

# Disable USB autosuspend for printers — prevents sleep-related latency
SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="07", ACTION=="add", \
    TEST=="power/autosuspend", ATTR{power/autosuspend}="-1"
SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="07", ACTION=="add", \
    TEST=="power/control", ATTR{power/control}="on"

# Regenerate Avahi AirPrint services after CUPS has had time to detect
# the new printer (run via a short delay script to avoid racing CUPS)
SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="07", ACTION=="add", \
    RUN+="/bin/sh -c 'sleep 3 && /opt/printserver/scripts/configure-avahi.sh &'"

# --- Printer disconnection: log only ---
SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="07", ACTION=="remove", \
    RUN+="/usr/bin/logger -t printer-hotplug 'USB printer disconnected'"
