# config/udev/99-printer.rules
# udev rules for USB printer hotplug
# /etc/udev/rules.d/99-printer.rules
#
# On connect:  hotplug-printer.sh add
#   → reloads CUPS (SIGHUP), waits for the scheduler, auto-adds the printer
#     queue, enables job acceptance, marks as shared, refreshes Avahi (mDNS).
#   → Windows 10/11, iOS and Android discover the printer automatically once
#     Avahi broadcasts the IPP/dnssd record.
#
# On disconnect: hotplug-printer.sh remove
#   → CUPS queue is preserved so queued jobs survive a cable reconnect.
#     Avahi is refreshed to remove stale mDNS advertisements.

# Reload (not restart!) CUPS so it re-enumerates USB devices without dropping
# active IPP connections from Windows clients.  A full restart kills all TCP
# connections, and Windows' IPP port monitor caches the dead state — the user
# then has to remove/re-add the printer to print again.
SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="07", ACTION=="add", \
    RUN+="/bin/systemctl reload cups.service"

# Disable USB autosuspend for printers — prevents sleep-related latency
SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="07", ACTION=="add", \
    TEST=="power/autosuspend", ATTR{power/autosuspend}="-1"
SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="07", ACTION=="add", \
    TEST=="power/control", ATTR{power/control}="on"

# Run the full hotplug pipeline asynchronously (needs a live CUPS scheduler).
SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="07", ACTION=="add", \
    RUN+="/bin/bash -c 'sleep 3 && /opt/printserver/scripts/hotplug-printer.sh add >> /var/log/printer-hotplug.log 2>&1 &'"

# USB printer disconnect — ID_USB_INTERFACES survives after the class attr is gone.
SUBSYSTEM=="usb", ENV{ID_USB_INTERFACES}=="*:070102:*", ACTION=="remove", \
    RUN+="/bin/bash -c 'sleep 5 && /opt/printserver/scripts/hotplug-printer.sh remove >> /var/log/printer-hotplug.log 2>&1 &'"

# Set read/write permissions on the raw USB device node for lp tools.
SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="07", MODE="0666"
