# config/udev/99-printer.rules
# udev rules for USB printer hotplug
# /etc/udev/rules.d/99-printer.rules
#
# On connect:  hotplug-printer.sh add
#   → restarts CUPS, waits for the scheduler, auto-adds the printer queue,
#     enables job acceptance, marks as shared, then refreshes Avahi (mDNS).
#   → Windows 10/11, iOS and Android discover the printer automatically once
#     Avahi broadcasts the IPP/dnssd record.
#
# On disconnect: hotplug-printer.sh remove
#   → CUPS queue is preserved so queued jobs survive a cable reconnect.
#     Avahi is refreshed to remove stale mDNS advertisements.

# Restart CUPS so it can re-enumerate the USB device tree.
# hotplug-printer.sh then waits for the scheduler before doing lpadmin work.
#
# IMPORTANT: use systemd-run --no-block instead of `systemctl restart`.
# udev RUN+= blocks the udev worker until the command exits; a synchronous
# `systemctl restart` waits for the whole D-Bus transaction and can stall
# the udev queue under load (cascading USB events, daemon-reload races).
# --no-block returns immediately and lets systemd schedule the restart.
SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="07", ACTION=="add", \
    RUN+="/bin/systemd-run --no-block --unit=cups-restart-on-hotplug systemctl restart cups.service"

# Disable USB autosuspend for printers — prevents sleep-related latency
SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="07", ACTION=="add", \
    TEST=="power/autosuspend", ATTR{power/autosuspend}="-1"
SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="07", ACTION=="add", \
    TEST=="power/control", ATTR{power/control}="on"

# Run the full hotplug pipeline asynchronously (needs a live CUPS scheduler).
SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="07", ACTION=="add", \
    RUN+="/bin/bash -c 'sleep 3 && /opt/printserver/scripts/hotplug-printer.sh add >> /var/log/printer-hotplug.log 2>&1 &'"

# USB printer disconnect — ID_USB_INTERFACES survives after the class attr is gone.
SUBSYSTEM=="usb", ENV{ID_USB_INTERFACES}=="*:070102:*", ACTION=="remove", \
    RUN+="/bin/bash -c 'sleep 5 && /opt/printserver/scripts/hotplug-printer.sh remove >> /var/log/printer-hotplug.log 2>&1 &'"

# Set read/write permissions on the raw USB device node for lp tools.
SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="07", MODE="0666"
