# config/systemd/printserver-web.service
[Unit]
Description=WiFi Print Server Web Interface
Documentation=https://github.com/your-repo/rpi-printserver
After=network-online.target cups.service avahi-daemon.service
Wants=network-online.target avahi-daemon.service
Requires=cups.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/opt/printserver
Environment=FLASK_ENV=production
Environment=PYTHONPATH=/opt/printserver
Environment=PRINTSERVER_CONFIG=/etc/printserver/config.ini

# Wait for CUPS to be ready before starting (up to 20s, then start anyway)
ExecStartPre=/bin/sh -c 'for i in 1 2 3 4 5 6 7 8 9 10; do lpstat -r 2>/dev/null && exit 0; sleep 2; done; echo "CUPS not ready after 20s, starting anyway"; exit 0'
ExecStart=/opt/printserver/venv/bin/gunicorn \
    --workers 1 \
    --threads 4 \
    --worker-class gthread \
    --bind 0.0.0.0:5000 \
    --timeout 30 \
    --graceful-timeout 10 \
    --keep-alive 5 \
    --max-requests 500 \
    --max-requests-jitter 50 \
    --worker-tmp-dir /dev/shm \
    --access-logfile - \
    --error-logfile - \
    'web.app:create_app()'
ExecReload=/bin/kill -HUP $MAINPID

# Restart configuration
Restart=always
RestartSec=5

# Timeouts
TimeoutStartSec=30
TimeoutStopSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/var/log/printserver
RuntimeDirectory=printserver

# Resource limits — tuned for 1GB Raspberry Pi
# MemoryHigh: soft limit — kernel applies memory pressure (reclaim/throttle)
# MemoryMax: hard limit — OOM kills the process to protect the rest of the OS
# Watchdog thread warns at 150 MB and clears caches at 190 MB RSS.
MemoryHigh=160M
MemoryMax=256M
MemorySwapMax=64M
CPUQuota=50%
OOMPolicy=kill

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=printserver-web

[Install]
WantedBy=multi-user.target
