# config/systemd/printer-wake.service
[Unit]
Description=Wake USB printer when a new print job arrives
After=cups.service

[Service]
Type=oneshot
ExecStart=/opt/printserver/scripts/wake-printer.sh
TimeoutStartSec=15
# Prevent rapid re-triggering if many jobs arrive at once
# (systemd will coalesce triggers within a short window)
