# config/systemd/printer-watchdog.timer
[Unit]
Description=Run printer health watchdog every 30 seconds

[Timer]
OnBootSec=30
OnUnitActiveSec=30
AccuracySec=5

[Install]
WantedBy=timers.target
