# config/systemd/printserver-web.service
[Unit]
Description=WiFi Print Server Web Interface
Documentation=https://github.com/your-repo/rpi-printserver
After=network.target cups.service avahi-daemon.service
Wants=cups.service avahi-daemon.service
Requires=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/printserver
Environment=FLASK_ENV=production
Environment=PYTHONPATH=/opt/printserver
Environment=PRINTSERVER_CONFIG=/etc/printserver/config.ini
ExecStart=/usr/bin/python3 -c "from web.app import create_app; app = create_app(); app.run(host='0.0.0.0', port=5000)"
ExecReload=/bin/kill -HUP $MAINPID

# Restart configuration
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/var/log/printserver

# Resource limits
MemoryMax=256M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=printserver-web

[Install]
WantedBy=multi-user.target
